# cargobike.service
[Unit]
Description=Cargobike

[Service]
ExecStart=/usr/local/bin/cargobike -version $(date +%s)

[Install]
WantedBy=multi-user.target
